# services/mpv-tts/run
#!/data/data/com.termux/files/usr/bin/sh
exec 2>&1
SOCK="$PREFIX/tmp/mpv-tts.sock"
rm -f "$SOCK"
exec mpv --idle=yes --force-window=no --input-ipc-server="$SOCK" --volume=100 --keep-open=no
